# docker/odoo.conf.conf/odoo.conf
[program:odoo]
command=/srv/odoo/odoo-bin -d odoopbx_%(ODOO_VERSION)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
