# docker/asterisk.conf
[program:asterisk]
command=/usr/sbin/asterisk -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
