# docker/supervisord.d/salt-minion.ini
[program:salt-minion]
command=/usr/local/bin/salt-minion
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=6
autostart=%(ENV_MINION_AUTOSTART)s
