# docker/supervisord.d/salt-master.ini
[program:salt-master]
command=/usr/local/bin/salt-master
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=4
autostart=%(ENV_MASTER_AUTOSTART)s
