# docker/supervisord.d/salt-api.ini
[program:salt-api]
command=/usr/local/bin/salt-api
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=5
autostart=%(ENV_API_AUTOSTART)s
